# laptop/com.user.claude-token-push.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heenal.claude-token-push</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/heenal/openclaw-claude-oauth-hybrid-sync/laptop/push-claude-tokens.sh</string>
    </array>
    <!-- Push every 4 hours (14400 seconds) -->
    <key>StartInterval</key>
    <integer>14400</integer>
    <!-- Also push on load (when Mac wakes up) -->
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/claude-token-push.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-token-push.err</string>
    <!-- Don't let it pile up -->
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
